# mental-wellness-tracker_1.0.0-1_all/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable systemd service
if [ -d /run/systemd/system ]; then
    systemctl stop mental-wellness-tracker.service || true
    systemctl disable mental-wellness-tracker.service || true
fi
